# installers/debian/poseidon.service
[Unit]
Description = Poseidon
After = docker.service

[Service]
ExecStart = /usr/bin/poseidon

[Install]
WantedBy = multi-user.target
