# installers/debian/prerm
#!/bin/sh

set -e

systemctl stop poseidon.service
systemctl disable poseidon.service
